# Minecraft/minecraft.service
[Unit]
Description=Minecraft

[Service]
WorkingDirectory=/minecraft
ExecStart=changeme
